# ops/systemd/utility-watershed-analytics-backup.timer
[Unit]
Description=Run the Utility Watershed Analytics database backup daily

[Timer]
OnCalendar=*-*-* 03:00:00 UTC
RandomizedDelaySec=15m
Persistent=true
Unit=utility-watershed-analytics-backup.service

[Install]
WantedBy=timers.target
