# ops/systemd/utility-watershed-analytics-backup-retention.timer
[Unit]
Description=Apply Utility Watershed Analytics backup retention weekly

[Timer]
OnCalendar=Sun *-*-* 04:00:00 UTC
RandomizedDelaySec=15m
Persistent=true
Unit=utility-watershed-analytics-backup-retention.service

[Install]
WantedBy=timers.target
