# docs/work-packages/20260716-db03-production-runtime-convergence/artifacts/wepp3-interim-runtime.service
[Unit]
Description=Utility Watershed Analytics interim production runtime
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target
ConditionPathExists=/etc/utility-watershed-analytics/runtime.env
ConditionPathExists=/etc/utility-watershed-analytics/database-identity
ConditionPathExists=/workdir/utility-watershed-analytics/compose.db03.yml

[Service]
Type=oneshot
RemainAfterExit=yes
User=root
Group=root
WorkingDirectory=/workdir/utility-watershed-analytics
ExecStartPre=/usr/local/lib/utility-watershed-analytics-runtime/db03-20260717-61e1667/scripts/check_runtime_environment.sh --env-file /etc/utility-watershed-analytics/runtime.env --expected-uid 0
ExecStart=/usr/local/lib/utility-watershed-analytics-runtime/db03-20260717-61e1667/scripts/start_runtime.sh --compose-file /workdir/utility-watershed-analytics/compose.db03.yml --project-name utility-watershed-analytics --env-file /etc/utility-watershed-analytics/runtime.env --identity-file /etc/utility-watershed-analytics/database-identity
ExecReload=/usr/local/lib/utility-watershed-analytics-runtime/db03-20260717-61e1667/scripts/with_operation_lock.sh --mode exclusive -- /usr/bin/docker compose --project-name utility-watershed-analytics --env-file /etc/utility-watershed-analytics/runtime.env --file /workdir/utility-watershed-analytics/compose.db03.yml up --detach --no-build --no-deps --pull never server caddy
ExecStop=/usr/local/lib/utility-watershed-analytics-runtime/db03-20260717-61e1667/scripts/with_operation_lock.sh --mode exclusive -- /usr/bin/docker compose --project-name utility-watershed-analytics --env-file /etc/utility-watershed-analytics/runtime.env --file /workdir/utility-watershed-analytics/compose.db03.yml stop --timeout 30 server caddy
TimeoutStartSec=15m
TimeoutStopSec=2m

[Install]
WantedBy=multi-user.target
