# ops/systemd/utility-watershed-analytics.service.d/30-db30b-artifact-route.conf
[Service]
ExecStartPost=/usr/local/lib/utility-watershed-analytics-runtime/db05-20260717-6f46aaf/scripts/with_operation_lock.sh --mode exclusive -- /usr/local/lib/utility-watershed-analytics-artifact-serving/ensure_artifact_proxy.sh
ExecReload=
ExecReload=/usr/local/lib/utility-watershed-analytics-runtime/db05-20260717-6f46aaf/scripts/with_operation_lock.sh --mode exclusive -- /usr/local/lib/utility-watershed-analytics-artifact-serving/ensure_artifact_proxy.sh
